# conf.d/xvfb.conf
[program:xvfb]
command=Xvfb :0 -screen 0 720x480x24
autorestart=true
