# conf.d/xvfb.conf
[program:xvfb]
command=Xvfb :0 -screen 0 2400x1080x24
autorestart=true
